# etc/sv/gerbera/run
#!/usr/bin/env /lib/runit/invoke-run

set -e

ping -q -w 1 -c 1 ff02::2 >/dev/null || exit 1

exec 2>&1
echo "Starting $DESC: $NAME."

LOGFILE=/var/log/$NAME.log

DAEMON_ARGS="-c $GERBERA_CONFIG -l $LOGFILE"

touch "$LOGFILE"
chmod 640 "$LOGFILE"
chown gerbera:adm "$LOGFILE"

exec chpst -u gerbera:gerbera /usr/bin/$NAME $DAEMON_ARGS
